# cron/crontab.TmEe0I/crontab
*/30 * * * * /home/x/Bureau/PROJET/job13/get-logs.sh
